# distro/deb/debian/python-fwmacro.postinst
#!/bin/sh
set -e

update-rc.d fwmacro start 39 S . stop 39 1 .
cat /etc/fwmacro/README
exit 0
